# devops/pythonsucks.socket
[Unit]
Description=pythonsucks socket

[Socket]
ListenStream=/run/pythonsucks.gunicorn.sock

[Install]
WantedBy=sockets.target
